# launchagents/Library/LaunchAgents/com.tnez.tmux-scratch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tnez.tmux-scratch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <string>/opt/homebrew/bin/tmux new-session -ds scratch -c /Users/tnez/PARA/DESK/inbox 2>/dev/null; /opt/homebrew/bin/tmux send-keys -t scratch 'nvim /Users/tnez/PARA/DESK/inbox/scratch.md' Enter</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tmux-scratch.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tmux-scratch.err</string>
</dict>
</plist>
